# com.example.gunicorn.plist
<!-- com.example.gunicorn.plist -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.gunicorn</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/gunicorn</string>
        <string>--config</string>
        <string>/tmp/django_app/web_app/gunicorn_config.py</string>
        <string>web_app.wsgi:application</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/tmp/django_app/web_app</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
